# config/certbot-weekly-renew.crontab
33 5 * * 5 certbot renew --post-hook "systemctl reload nginx"
